# systemd/code-router.service
[Unit]
Description=Refresh OAuth token for claude-code-router (icode)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# CCR's daemon (spawned by `ccr restart` inside the refresh script) inherits
# this env, so the corporate CA bundle MUST be exported here -- otherwise
# the freshly spawned daemon cannot validate TLS to the upstream gateway.
Environment=NODE_EXTRA_CA_CERTS=%h/.local/share/ca-certs/code-router-ca.pem
ExecStart=%h/bin/code-router-refresh-token
StandardOutput=journal
StandardError=journal
